# tomboy@.service
[Unit]
Description=Tomboy
After=xorg.target

[Service]
ExecStart=/usr/bin/tomboy
Environment=DISPLAY=%i

[Install]
WantedBy=startup.target
