# xcape@.service
[Unit]
Description=Run xcape
After=xorg.target

[Service]
ExecStart=/usr/bin/xcape
Environment=DISPLAY=%i

[Install]
WantedBy=startup.target
